# crontab
# @reboot sleep 30 && python /home/gwlillie/Thonny/blinkingLED.py &
# @reboot sleep 60 && python /home/gwlillie/Thonny/start.py &
@reboot sleep 60 && python /home/gwlillie/Thonny/buttonStream.py &
